# setup_res/etc/cron.hourly/shutter
#!/bin/bash
#UPT=$(</proc/uptime)
#UPT=${UPT%%.*}
#[ $UPT -gt 880000 ] && poweroff

